# crates/dormantd/systemd/dormant.service
[Unit]
Description=dormant — proximity-driven display blanking daemon
Documentation=https://github.com/legion-works/dormant
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Cap glibc malloc arenas so freed screensaver/render heap is returned to the OS
# instead of retained per-thread (pairs with the daemon's worker_threads=2 and the
# malloc_trim after screensaver teardown — see crates/dormantd/src/main.rs).
Environment=MALLOC_ARENA_MAX=2
ExecStart=%h/.local/bin/dormantd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
